# systemd-services/root/backrest-backup/backrest-backup.timer
[Unit]
Description=Run backrest data backup daily at 00:10

[Timer]
OnCalendar=*-*-* 00:00:10
Persistent=true

[Install]
WantedBy=timers.target
