# systemd-services/root/backrest-backup/backrest-backup.service
[Unit]
Description=Backup Backrest Data
After=network.target

[Service]
Type=oneshot
ExecStart=/home/dpi0/scripts/backrest-backup.sh
User=root
Group=root
